# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:chatbot]
command=python /app/chatbot_service_sse.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/chatbot_err.log
stdout_logfile=/var/log/supervisor/chatbot_out.log
environment=CHATBOT_SERVICE_PORT=5001,CHATBOT_SERVICE_HOST=0.0.0.0,PYTHONUNBUFFERED=1
user=root

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nginx_err.log
stdout_logfile=/var/log/supervisor/nginx_out.log
user=root

